# contrib/supervisord.conf.in
[supervisord]
nodaemon=true

[program:{DAEMON}d]
command=/srv/explorer/{DAEMON}/bin/{DAEMON}d -conf=/data/.{DAEMON}.conf -datadir=/data/{DAEMON}
directory=/data

[program:bitcoind_for_liquid]
command=/srv/explorer/bitcoin/bin/bitcoind -conf=/data/.bitcoin.conf -datadir=/data/bitcoin
directory=/data
autostart={ISLIQUID}

[program:electrs]
command=/srv/explorer/electrs_{DAEMON}/bin/electrs {ELECTRS_DEBUG} --timestamp --http-addr 127.0.0.1:3000 --network {ELECTRS_NETWORK} {PARENT_NETWORK} --daemon-dir /data/{DAEMON} --monitoring-addr 0.0.0.0:4224 --db-dir /data/electrs_{DAEMON}_db/{NETWORK} --precache-scripts /srv/explorer/popular-scripts.txt
autostart={EXPLORERAUTOSTART}
directory=/data
stopasgroup=true
stderr_logfile=/data/logs/electrs-{DAEMON}.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=5

[program:prerenderer]
environment=STATIC_ROOT=/{NGINX_PATH},BASE_HREF=/{NGINX_PATH}nojs/,NODE_ENV=production,API_URL=http://localhost/{NGINX_PATH}api,TZ=UTC
command=/bin/bash -c 'source /root/.nvm/nvm.sh && cd /srv/explorer/source && npm run prerender-server -- --cluster {FLAVOR}'
autostart={EXPLORERAUTOSTART}
directory=/srv/explorer/source
stopasgroup=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart={EXPLORERAUTOSTART}

[program:tor]
command=/usr/bin/tor
